# debian/demoapi.service
[Unit]
Description=Demo Api Service
After=network.target

[Service]
Type=simple
ExecStart=/opt/bluekiri/demoapi/demoapi
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
